# misc/pcilib_html.service.in
[Unit]
Description=Pcilib Python html server
Wants=network.target nss-lookup.target
After=network.target nss-lookup.target
Requires=pcilib_api.service
After=pcilib_api.service

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/pcilib
ExecStart=${PYTHON_EXECUTABLE} ${PCILIB_SERVER_DIR}/pcilib_html_server.py -d $SERVER_DEVICE -m $SERVER_MODEL -p $HTML_SERVER_PORT -e --api-server-port $API_SERVER_PORT

[Install]
WantedBy=multi-user.target
